# distro/debian/powertweakd.conffiles
/etc/powertweak/tweaks
/etc/init.d/powertweakd
/etc/default/powertweakd
